# supervisor.d/mssqlserver.conf
# MS SQL Server
[program:sqlserver]
command=/opt/mssql/bin/sqlservr
stderr_logfile = /var/log/supervisord/sqlserver-stderr.log
stdout_logfile = /var/log/supervisord/sqlserver-stdout.log

