# supervisor.d/samba.conf
# Samba
[program:samba]
command=/usr/sbin/smbd -F
stderr_logfile = /var/log/supervisor/samba-stderr.log
stdout_logfile = /var/log/supervisor/samba-stdout.log

